# cfg/upload.crontab
12 */4 * * * bash -c "/app/upload.sh >> /app/cron.log 2>&1"
# Run upload
